# supervisord.conf
; supervisord.conf — manages Nginx and PHP-FPM as two processes
; inside one Docker container.
; supervisord starts both, restarts either if they crash,
; and forwards their logs to stdout so Render can see them.

[supervisord]
nodaemon=true              ; run in foreground (required for Docker)
logfile=/dev/null          ; log to stdout instead of a file
logfile_maxbytes=0

; ── PHP-FPM ──────────────────────────────────────────────────────
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
autostart=true             ; start when supervisord starts
autorestart=true           ; restart if it crashes
stdout_logfile=/dev/stdout ; forward logs to Docker stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Nginx ─────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
